# edgevpn/assets/edgevpn.openrc
#!/sbin/openrc-run

description="EdgeVPN service"
output_log="/var/log/edgevpn.log"
error_log="/var/log/edgevpn.err"
pidfile=/var/run/edgevpn.pid
command="/usr/local/bin/edgevpn"
command_args="--api --dhcp"
command_background="yes"

# Load configuration from /etc/conf.d/edgevpn
if [ -f /etc/conf.d/edgevpn ]; then
    . /etc/conf.d/edgevpn
    # Export environment variables for the command
    export EDGEVPNTOKEN
fi

depend() {
    need net
    after firewall
} 